# crontab
# speedtest log
5 * * * * $HOME/speedtest.sh